# etc/launchd/daemons/org.openzfsonosx.zed.service.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.openzfsonosx.zed.service</string>
	<key>ProgramArguments</key>
	<array>
		<string>@launchdscriptdir@/zed.service.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/private/var/log/openzfsonosx.err</string>
	<key>StandardOutPath</key>
	<string>/private/var/log/openzfsonosx.log</string>
	<key>WatchPaths</key>
	<array>
		<string>/dev/zfs</string>
		<string>/etc/zfs</string>
	</array>
</dict>
</plist>
